# emocute-toolkit/automation/launchd/com.emocute.toolkit.weekly-downloads-clean.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 週次 Downloads 直下クリーンアップ提案 (dry-run のみ、実移動は手動 --apply)。日曜 09:00。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.emocute.toolkit.weekly-downloads-clean</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/bin/emocute</string>
    <string>ops</string>
    <string>downloads-weekly-clean</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/automation/logs/weekly-downloads-clean.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/automation/logs/weekly-downloads-clean.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
