# emocute-toolkit/automation/launchd/com.emocute.toolkit.daily-cookie-expiry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 日次 Playwright cookie 期限監視 (14 日以内期限切れを warn)。08:30。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.emocute.toolkit.daily-cookie-expiry</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/bin/emocute</string>
    <string>infra</string>
    <string>cookie-expiry-monitor</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/automation/logs/daily-cookie-expiry.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/automation/logs/daily-cookie-expiry.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
