# emocute-toolkit/automation/launchd/com.emocute.toolkit.weekly-failure-mode-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 週次 Sale 配下の failure-mode 早期検出。火曜 08:00。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.emocute.toolkit.weekly-failure-mode-check</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/bin/emocute</string>
    <string>ops</string>
    <string>failure-mode-early-check</string>
    <string>/Users/emocute/Downloads/Sale</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>2</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/automation/logs/weekly-failure-mode-check.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/emocute/Downloads/Sketches/emocute-toolkit/automation/logs/weekly-failure-mode-check.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
